# etc/asynk/supervisor.conf
[program:asynk]
command=/path/to/project/bin/comma AsynkWorker
process_name=%(program_name)s_%(process_num)04d
numprocs=5
autostart=true
autorestart=true

